# crontab
*/1 * * * * cd /opt/ddns/ && python ddns.py >> /tmp/ddns.log 2>&1
