# crontab
*/1 * * * * cd /web/ddns/ && python ddns.py >> ddns.log >>2&1
