# systemd/ddns.timer
[Unit]
Description=Runs ddns every minute

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
